# build/celery_beat_docker.conf
; ================================
;  celery beat supervisor
; ================================

[program:celerybeat]
directory=/opt/webtetrado
command=celery -A WebTetrado beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler

user=root
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
